# systemd/sonatum-stack.service
[Unit]
Description=Sonatum Music Docker Compose stack
Requires=docker.service
After=docker.service network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/sonatum
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
